# fs/etc/supervisor/conf.d/30-nginx.conf
[program:nginx]
command=nginx -g 'pid /run/nginx/nginx.pid; daemon off;'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
